# screenpulse-analyzer.service
[Unit]
Description=ScreenPulse Analyzer - AI Video Analysis Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/user/screenpulse
ExecStart=/home/user/screenpulse/.venv/bin/python /home/user/screenpulse/screenpulse.py
Restart=always
RestartSec=10
StandardOutput=append:/home/user/screenpulse/screenpulse.log
StandardError=append:/home/user/screenpulse/screenpulse.log

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/user/screenpulse/.venv/bin:/usr/bin:/bin"
Environment="GEMINI_API_KEY=AIzaSyBA8ZU0IJvEqqIrJnzLmtF5yjxBGvBWXcw"

[Install]
WantedBy=default.target
