# screenpulse.service
[Unit]
Description=ScreenPulse - Motion-triggered screen recorder
After=graphical-session.target

[Service]
Type=forking
User=%u
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
WorkingDirectory=/home/user/.local/share/screenpulse
ExecStart=/home/user/.local/share/screenpulse/screenpulse-ctl start
ExecStop=/home/user/.local/share/screenpulse/screenpulse-ctl stop
PIDFile=/home/user/.local/share/screenpulse/screenpulse.pid
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
